# crontab
* * * * * root python3 /app/run_parsers.py >> /var/log/cron.log 2>&1